# openwrt/files/etc/rc.local
/etc/init.d/komari-agent start
/usr/sbin/smartdns -c /etc/smartdns/smartdns.conf

exit 0